# cron/crontab

0 * * * * /app/crontask.sh >> /var/log/cron.log 2>&1
# TODO: aanpassen naar elk uur